# com.namsung.michunggu.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.namsung.michunggu</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jeongwoo/code/dnjfakrkaalcjdrn/run.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jeongwoo/code/dnjfakrkaalcjdrn</string>

    <!-- 매월 1~7일 00:10 실행 (공휴일/주말이면 첫 영업일까지 대기) -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Day</key><integer>1</integer><key>Hour</key><integer>0</integer><key>Minute</key><integer>10</integer></dict>
        <dict><key>Day</key><integer>2</integer><key>Hour</key><integer>0</integer><key>Minute</key><integer>10</integer></dict>
        <dict><key>Day</key><integer>3</integer><key>Hour</key><integer>0</integer><key>Minute</key><integer>10</integer></dict>
        <dict><key>Day</key><integer>4</integer><key>Hour</key><integer>0</integer><key>Minute</key><integer>10</integer></dict>
        <dict><key>Day</key><integer>5</integer><key>Hour</key><integer>0</integer><key>Minute</key><integer>10</integer></dict>
        <dict><key>Day</key><integer>6</integer><key>Hour</key><integer>0</integer><key>Minute</key><integer>10</integer></dict>
        <dict><key>Day</key><integer>7</integer><key>Hour</key><integer>0</integer><key>Minute</key><integer>10</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/jeongwoo/code/dnjfakrkaalcjdrn/logs/launchd-out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jeongwoo/code/dnjfakrkaalcjdrn/logs/launchd-err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/jeongwoo</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
